# dummy.sv.conf
[program:dummy]
command=/bin/bash -c 'while true; do echo "$FOO - $(date)" > /tmp/date; sleep 1; done'
redirect_stderr=true
